# tap-station.service
[Unit]
Description=NFC Tap Logger Station
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/nfc-tap-logger
Environment="PATH=/home/pi/nfc-tap-logger/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pi/nfc-tap-logger/venv/bin/python -m tap_station.main
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/nfc-tap-logger/logs/tap-station.log
StandardError=append:/home/pi/nfc-tap-logger/logs/tap-station-error.log

# Restart on failure
StartLimitInterval=200
StartLimitBurst=5

# Run with elevated privileges for GPIO/I2C access
# Alternative: Add user to gpio and i2c groups
# SupplementaryGroups=gpio i2c

[Install]
WantedBy=multi-user.target
