# tap-station.service
# FlowState Station Service - TEMPLATE
#
# This is a template file. The install.sh script automatically generates
# the correct service file with your system's paths.
#
# If manually installing, replace these placeholders:
#   - YOUR_USERNAME: Your system username (e.g., pi, dc-nfc-1)
#   - /home/YOUR_USERNAME/flowstate: Path to your FlowState installation
#
# Manual installation:
#   1. Edit this file with your values
#   2. Copy to: sudo cp tap-station.service /etc/systemd/system/
#   3. Enable: sudo systemctl enable tap-station
#   4. Start: sudo systemctl start tap-station

[Unit]
Description=FlowState Station Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/flowstate
Environment="PATH=/home/YOUR_USERNAME/flowstate/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/YOUR_USERNAME/flowstate/venv/bin/python -m tap_station.main --config /home/YOUR_USERNAME/flowstate/config.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Restart on failure
StartLimitInterval=200
StartLimitBurst=5

# Run with elevated privileges for GPIO/I2C access
# Alternative: Add user to gpio and i2c groups
# SupplementaryGroups=gpio i2c

[Install]
WantedBy=multi-user.target
